# root/etc/s6-overlay/s6-rc.d/service-goaccess/run
#!/command/with-contenv bash

umask "${UMASK}"

exec s6-setuidgid goaccess /usr/local/bin/goaccess.sh
